# rootfs/etc/s6-overlay/s6-rc.d/init-usermap/up
with-contenv /usr/local/bin/init-usermap.sh
