# rootfs/etc/s6-overlay/s6-rc.d/git-watcher/finish
#!/command/with-contenv sh
# Final commit on service shutdown — catch anything in the debounce window.
# The main script also has a SIGTERM trap, but this is a safety net.
# Run as the obsidian user to match file ownership.
#
# BRANCH is expanded by the outer shell (where with-contenv guarantees
# the container environment is loaded).
BRANCH="${OBSIDIAN_GIT_BRANCH:-main}"
su-exec obsidian sh -c "
  cd /vault 2>/dev/null || exit 0
  git add -A 2>/dev/null || exit 0
  git diff --cached --quiet 2>/dev/null && exit 0
  git commit -m \"shutdown sync: \$(date -u '+%Y-%m-%d %H:%M:%S UTC')\" 2>/dev/null || exit 0
  git push origin \"${BRANCH}\" 2>/dev/null || true
"
