# rootfs/etc/s6-overlay/s6-rc.d/git-watcher/run
#!/command/with-contenv sh
# Git filesystem watcher service
# Drop privileges to the obsidian user (UID/GID set by init-usermap)
exec su-exec obsidian /usr/local/bin/git-watcher.sh
