# rootfs/etc/s6-overlay/s6-rc.d/git-watcher/run
#!/bin/sh
# Git filesystem watcher service
exec /usr/local/bin/git-watcher.sh
