# rootfs/etc/s6-overlay/s6-rc.d/init-config/up
with-contenv /usr/local/bin/init-config.sh
