# rootfs/etc/s6-overlay/s6-rc.d/obsidian-sync/run
#!/command/with-contenv sh
# Obsidian headless continuous sync service
# Drop privileges to the obsidian user (UID/GID set by init-usermap)
#
# Prefix every line with [obsidian-sync] to match the git-watcher log format.
# This makes interleaved Docker logs easy to read.
#
# sed -u (unbuffered) ensures lines appear immediately.
# When ob exits, the pipe breaks, sed gets EOF and exits, s6 restarts all.
# On shutdown, s6 sends SIGTERM to the service's process group — both ob
# and sed receive it.
# Note: no 'exec' — it has no effect before a pipeline in POSIX sh.
# The shell stays as s6's direct child and manages the pipeline process group.
su-exec obsidian sh -c 'exec ob sync --continuous --path /vault 2>&1' \
  | sed -u 's/^/[obsidian-sync] /'
