# rootfs/etc/s6-overlay/s6-rc.d/obsidian-sync/run
#!/command/with-contenv sh
# Obsidian headless continuous sync service
# Drop privileges to the obsidian user (UID/GID set by init-usermap)
#
# Uses log-prefix.sh to add [obsidian-sync] to every output line, matching
# the [git-watcher] convention. The wrapper avoids shell pipeline problems:
#   - Signals: SIGTERM reaches ob directly (not lost to sed)
#   - Exit code: ob's real exit code is returned to s6
#   - No SIGPIPE: ob writes to a FIFO, not a pipe to sed
exec su-exec obsidian /usr/local/bin/log-prefix.sh "[obsidian-sync]" \
  ob sync --continuous --path /vault
