# rootfs/etc/s6-overlay/s6-rc.d/obsidian-sync/run
#!/bin/sh
# Obsidian headless continuous sync service
exec ob sync --continuous --path /vault
