# rootfs/etc/s6-overlay/s6-rc.d/obsidian-sync/run
#!/bin/sh
# Obsidian headless continuous sync service
# Drop privileges to the obsidian user (UID/GID set by init-usermap)
exec su-exec obsidian ob sync --continuous --path /vault
